# deploy/supervisord.conf
; 单容器内拉起三进程：声纹(8101) + 上下文(8102) + Go实时核心(8000)。
; 各进程日志直接转到容器 stdout/stderr；崩溃自动重启。
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:voiceprint]
command=uvicorn app:app --host 127.0.0.1 --port 8101
directory=/app/voiceprint
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:context]
command=uvicorn app:app --host 127.0.0.1 --port 8102
directory=/app/context
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:backend]
command=/app/interview-backend
directory=/app
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
